# supervisor/scheduler.conf
[program:scheduler]
command=/var/www/bin/scheduler.sh
numprocs=1
autostart=true
autorestart=true
startretries=2
